# deploy-ova-with-vAPP-properties/netconfig.service
[Unit]
Description=Configure Network with vAPP properties
 
[Service]
ExecStart=/usr/local/bin/netconfig.sh
 
[Install]
WantedBy=multi-user.target
